# specific/gambit.crontab
0 5 * * * cd ~/Projects/inspiration && ./update.sh
